# etc/rc.local
sh /home/pi/scripts/start_irtoy.sh
perl /home/pi/scripts/read_meter.pl &
